# packaging/herdr-gateway-dev.service
[Unit]
Description=herdr-gateway (dev) — runs the repo build output as the live instance
Documentation=https://github.com/vantt/herdr-gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Dev-as-live: ExecStart points straight at the repo's release binary, and
# WorkingDirectory is the repo, so `static/` and the zero-config default resolve
# in place. Rebuild + `systemctl --user restart herdr-gateway-dev` = the new
# build is live — no copy, no reinstall (see dev-deploy.sh).
WorkingDirectory=@REPO@
ExecStart=@REPO@/target/release/herdctl
# Optional secrets (GitHub/Telegram); the web login token is auto-managed by the
# binary in ~/.config/herdr-gateway/herdctl.env.
EnvironmentFile=-@ENVFILE@
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
