# packaging/herdr-go.service
[Unit]
Description=herdr-go — web remote gateway + supervisor for herdr
Documentation=https://github.com/vantt/herdr-go
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@BIN@ --config @CONFIG@
EnvironmentFile=@ENVFILE@
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=@DATA_DIR@ @CONFIG_DIR@

[Install]
WantedBy=default.target
