# packaging/herdr-go-dev.service
[Unit]
Description=herdr-go (dev) — runs the repo build output as the live instance
Documentation=https://github.com/vantt/herdr-go
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO@
ExecStart=@REPO@/target/release/herdr-go
EnvironmentFile=-@ENVFILE@
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=@DATA_DIR@ @CONFIG_DIR@ @REPO@
ProtectHome=read-only

[Install]
WantedBy=default.target
