# packaging/herdr-gateway.service
[Unit]
Description=herdr-gateway (herdctl) — web remote gateway + supervisor for herdr
Documentation=https://github.com/thanhsmind/herdr-gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# One self-healing unit: the gateway supervises herdr, systemd supervises the
# gateway. Restart on any exit; back off so a crash loop does not hammer.
ExecStart=@BIN@ --config @CONFIG@
EnvironmentFile=@ENVFILE@
Restart=always
RestartSec=3
# Light hardening — the gateway still runs as the operator (needs herdr access).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/.local/share/herdr-gateway %h/.config/herdr-gateway
ProtectHome=read-only

[Install]
WantedBy=default.target
